# collections/ansible_collections/jctanner/k8s_lab/roles/dns_server/files/godns.service
[Unit]
Description=Custom DNS Server
After=network.target

[Service]
ExecStart=/opt/godns/godns.sh

[Install]
WantedBy=multi-user.target
